# Code - OSS/User/History/6aba7b38/2jbS.service
[Unit]
Description=Battery low alert

[Service]
ExecStart=%h/.config/scripts/battery_alert.sh
