# Code - OSS/User/History/6aba7b38/g20f.service
[Unit]
Description=Battery low alert

[Service]
ExecStart=%h/.local/bin/battery_alert.sh
